# devops/cron/crontab
# Laravel scheduler. Install for the app user on the scheduler host ONLY
# (exactly one host — running it on two would double every scheduled task):
#     crontab -u aiku devops/cron/crontab
* * * * * cd /home/aiku/aiku/current && /usr/bin/php8.4 artisan schedule:run >> /dev/null 2>&1
